# vaked-agents/ci/swe-af-orchestrator/deploy/swe-af.slice
# Aggregate cgroup fence for the swe-af batch on a worker host (e.g. bench-node).
# swe-af is a GUEST on the box: this slice bounds the orchestrator + all worker
# subprocesses so they can never starve the host's other services.
[Unit]
Description=swe-af batch slice (bounded guest)

[Slice]
MemoryMax=8G
MemoryHigh=6G
CPUQuota=600%
TasksMax=4096
