# daemons/synapsed/vaked-synapsed.service
[Unit]
Description=Vaked Synapsed P2P Capability-Graph Daemon Cluster Broker
After=network.target
StartLimitIntervalSec=0

[Service]
Type=notify
ExecStart=/usr/bin/env vaked-synapsed --config /etc/vaked/mesh-policy.vaked
WatchdogSec=5
Restart=always

# Hard System Isolation (No Root Leaks)
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_ADMIN
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK
DeviceAllow=/dev/null rwm
MemoryMax=512M
