# daemons/sandboxd/zig-out/bin/sandboxd.service
[Unit]
Description=sandboxd — namespace/exec enforcement daemon
After=network.target
[Service]
Type=simple
ExecStart=/usr/local/bin/sandboxd --policy /etc/sandboxd/policy.json
Restart=on-failure
RestartSec=5
[Install]
WantedBy=multi-user.target
